# etc/cron.d/datatracker
# -*- indent-with-tabs: 0 -*-
#
# This file controls the datatracker-related cronjobs.  It's intended to be
# symlinked into /etc/init.d/ .  It is a crontab in order to be able to add
# actions which don't fit into the daily or hourly pattern.  It, and daily
# and hourly scripts it triggers, are part of the datatracker repository in
# order to handle cronjob changes that go together with code changes well.

PATH=/sbin:/usr/sbin:/bin:/usr/bin

DTDIR=/a/www/ietf-datatracker/web

# m h dom mon dow user  command

05 00   *  *  *   wwwrun  test -x $DTDIR/bin/daily  && $DTDIR/bin/daily

05 1-23 *  *  *   wwwrun  test -x $DTDIR/bin/hourly && $DTDIR/bin/hourly


